# tools/beetle-updater/beetle-demo-service.service
[Unit]
Description=example systemd service unit file.

[Service]
ExecStart=/bin/bash -c "echo hi"

[Install]
WantedBy=multi-user.target
